# build/rpm/crispy.service
[Unit]
After=network.target

[Service]
Type=simple

EnvironmentFile=/etc/default/crispy
ExecStart=/usr/bin/crispy

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
